# thingino-build-broker.service
[Unit]
Description=Thingino web-builder broker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/thingino-broker
# secret.env holds the secrets only (chmod 600, root-owned):
#   GITHUB_TOKEN=github_pat_...        (fine-grained, scoped to the builder repo)
#   GITHUB_REPO=thingino/thingino-web-builder
EnvironmentFile=/opt/thingino-broker/secret.env
Environment=BIND_ADDR=127.0.0.1:8080
Environment=STATIC_DIR=/opt/thingino-broker/web
Environment=DEFCONFIGS_PATH=/opt/thingino-broker/defconfigs.json
Environment=DB_PATH=/opt/thingino-broker/broker.db
Environment=ROLLING_TAG=web-builds
Environment=GLOBAL_HOURLY_LIMIT=10
Environment=PER_IP_HOURLY_LIMIT=2
# Behind nginx/Caddy/Cloudflare, set this so the per-IP limit sees the real
# client and not the proxy (X-Forwarded-For or CF-Connecting-IP):
#Environment=IP_HEADER=X-Forwarded-For
ExecStart=/opt/thingino-broker/thingino-build-broker
Restart=on-failure
RestartSec=3

# Hardening
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/thingino-broker

[Install]
WantedBy=multi-user.target
